# example.ini
[program:example]
directory=/var/tmp
command=top
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/example.err.log
stdout_logfile=/var/log/supervisor/example.out.log